# docker/launchd/ai.agentdash.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.agentdash.agent</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/agentdash/bin/agentdash-compose-supervisor.sh</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>WorkingDirectory</key>
	<string>/opt/agentdash</string>

	<key>StandardOutPath</key>
	<string>/opt/agentdash/logs/launchd.out.log</string>

	<key>StandardErrorPath</key>
	<string>/opt/agentdash/logs/launchd.err.log</string>

	<key>ProcessType</key>
	<string>Background</string>

	<key>ThrottleInterval</key>
	<integer>10</integer>

	<key>ExitTimeOut</key>
	<integer>30</integer>
</dict>
</plist>
